# config/crontab
0 9-23 * * * php /home/app/src/bin/console app:process-cars >/dev/null 2>&1
